# releng/root-image/etc/rc.local
#!/bin/bash

echo Changing root shell to zsh
chsh -s /bin/zsh root

