# releng/root-image/etc/rc.local
#!/bin/bash

useradd -m -s /bin/bash -U -G wheel jreese

